# scripts/riju.service
[Unit]
Description=Riju online coding sandbox

[Service]
ExecStart=riju-serve

[Install]
WantedBy=multi-user.target
