# scripts/riju.service
[Unit]
Description=Riju online coding sandbox

[Service]
ExecStart=docker run -p 0.0.0.0:6119 riju

[Install]
WantedBy=multi-user.target
